# proxy-manager/rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Nginx Proxy Manager
# Runs the MySQL database server
# ==============================================================================

declare -a options

options+=(--"datadir=/data/mysql")
options+=(--"tmpdir=/tmp/")

exec /usr/bin/mysqld_safe "${options[@]}"
