# proxy-manager/rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Nginx Proxy Manager
# Runs the MySQL database server
# ==============================================================================
declare -a options

bashio::log.info "Starting MySQL database server..."

options+=(--datadir /data/mysql)
options+=(--tmpdir /tmp/)

exec s6-setuidgid mysql /usr/bin/mysqld "${options[@]}"
