# proxy-manager/rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Nginx Proxy Manager
# Runs the MySQL database server
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

declare -a options

options+=(--datadir /data/mysql)
options+=(--tmpdir /tmp/)

exec s6-setuidgid mysql /usr/bin/mysqld "${options[@]}"
