# proxy-manager/rootfs/etc/services.d/manager/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Nginx Proxy Manager
# Runs the nginx proxy manager interface
# ==============================================================================
declare -a options

# Wait for MySQL to become available
bashio::net.wait_for 3306

bashio::log.info "Starting the Manager..."

options+=(--abort_on_uncaught_exception)
options+=("--max_old_space_size=250")

export NODE_ENV=production

cd /opt/nginx-proxy-manager \
    || bashio::exit.nok "Could not change directory to app"

exec node "${options[@]}" /opt/nginx-proxy-manager/src/backend/index.js
