# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/npm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Nginx Proxy Manager
# Runs the nginx proxy manager interface
# ==============================================================================
declare -a options

bashio::log.info "Starting the Manager..."

options+=(--abort_on_uncaught_exception)
options+=("--max_old_space_size=250")

export NODE_ENV=production

if bashio::debug; then
  export DEBUG=1
fi


cd /opt/nginx-proxy-manager \
    || bashio::exit.nok "Could not change directory to app"

exec node "${options[@]}" /opt/nginx-proxy-manager/index.js
