# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/init-npm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Nginx Proxy Manager
# This file applies patches so the add-on becomes compatible
# ==============================================================================
mkdir -p \
    /config/access \
    /config/custom_ssl \
    /config/letsencrypt \
    /config/letsencrypt-acme-challenge \
    /config/logs \
    /config/nginx/dead_host \
    /config/nginx/default_host \
    /config/nginx/default_www \
    /config/nginx/proxy_host \
    /config/nginx/redirection_host \
    /config/nginx/stream \
    /config/nginx/temp \
    /tmp/letsencrypt-lib \
    /tmp/letsencrypt-log \
    /tmp/nginx/body \
    /tmp/nginx/cache/private\
    /tmp/nginx/cache/public \
    /tmp/nginx/proxy \
    /var/lib/nginx/

ln -s /tmp/nginx /var/tmp/nginx
ln -s /tmp/nginx /var/lib/nginx/tmp
ln -s /config /opt/nginx-proxy-manager/config

ln -sf /config/letsencrypt /etc/letsencrypt

# NGinx needs this file to be able to start.
# It will not continously log into it.
mkdir -p /var/lib/nginx/logs
touch /var/lib/nginx/logs/error.log
chmod 777 /var/lib/nginx/logs/error.log

# This file generates a dummy SSL certificate
if ! bashio::fs.directory_exists "/config/nginx/dummycert.pem" && ! bashio::fs.file_exists "/config/nginx/dummykey.pem";
then
  bashio::log.info "Generating dummy SSL certificate"
  openssl req \
    -new \
    -newkey rsa:2048 \
    -days 3650 \
    -nodes \
    -x509 \
    -subj '/O=Nginx Proxy Manager/OU=Dummy Certificate/CN=localhost' \
    -keyout /config/nginx/dummykey.pem \
    -out /config/nginx/dummycert.pem \
    || bashio::exit.nok "Could not generate dummy certificate"
fi
