# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Nginx Proxy Manager
# Runs the nginx proxy manager interface
# ==============================================================================

bashio::log.info "Starting NGinx..."

# Create required folders
mkdir -p \
  /tmp/nginx/body \
  /var/log/nginx \
  /data/custom_ssl \
  /data/logs \
  /data/access \
  /tmp/nginx/cache/public \
  /tmp/nginx/cache/private \
  /var/cache/nginx/proxy_temp

# Run
exec nginx
