# proxy-manager/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Nginx Proxy Manager
# Runs the DNSMasq
# ==============================================================================

declare -a options

# Allow dnsmasq to bind on ports < 1024
#setcap CAP_NET_BIND_SERVICE=+eip "$(command -v dnsmasq)"

options+=(--keep-in-foreground)

exec dnsmasq "${options[@]}"
