# proxy-manager/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Nginx Proxy Manager
# Runs the DNSMasq
# ==============================================================================
declare -a options

bashio::log.info "Starting DNSMasq..."

options+=(--keep-in-foreground)

exec dnsmasq "${options[@]}"
